# examples/systemd/dummy.env.service
# WRAP pattern (env variant): for programs that read the config path from an
# environment variable rather than a CLI flag. The `.secretrc` entry sets
# `env = "APP_CONFIG"`, so ss-fetcher exports APP_CONFIG=<materialised path>
# for the child. No {config} placeholder is needed in ExecStart.
#
# Install to: ~/.config/systemd/user/dummy.service

[Unit]
Description=Dummy service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/dummy

ExecStart=ss-fetcher run -- /opt/dummy/venv/bin/python3 /opt/dummy/app.py

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
